# etc/supervisord.d/hipache.conf
[program:hipache]
command=/usr/bin/hipache --config /etc/hipache.json
autorestart=true
stdout_logfile=/var/log/supervisor/hipache.log
redirect_stderr=true
